# script/sora-archive-uploader.timer
[Unit]
Description=Sora Archive Uploader Wake Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target